# projects/26.Askirom.VPS-CoreOS-Migration/quadlets/cloudflared.container
[Unit]
Description=Cloudflare Tunnel
After=network-online.target

[Container]
Image=docker.io/cloudflare/cloudflared:latest
Exec=tunnel run --token ${CLOUDFLARED_TOKEN}
EnvironmentFile=%h/.config/cloudflared/env
Network=host

[Service]
Restart=always

[Install]
WantedBy=default.target
